# sys/timers/docker-system-prune.service
[Unit]
Description=Prune the docker system

[Service]
Type=oneshot
ExecStart=docker image prune -f --filter "until=24h"
